# resources/bash/profile.d/welcome.sh


echo::info "User: $USER, UID: $USER_UID, GID: $(id -g), Groups: $USER_GROUPS"
echo::info ""
echo::info "dbctl is our server control script (database + control)"
echo::info ""
echo::info "Execute dbctl for a quick and easy server management"
echo::info ""
echo::info "  dbctl"
echo::info ""