# resources/bash/profile/profile.d/postgres-restore.sh


# color in diagnostic messages (https://www.postgresql.org/docs/current/app-pgrestore.html)
# value may be always, auto and never
export PG_COLOR=always
